# script/play.cron
* * * * * /Users/rupert/.n/bin/node /Users/rupert/Code/game-media/script/play.js
